# profile.d/bash_setup.sh
# vi: set ft=sh :

SHELL_IS_BASH=1

source ~/.profile.d/setup.sh
