# profile.d/zshell_setup.sh
# vi: set ft=zsh :

SHELL_IS_ZSHELL=1

source ~/.profile.d/setup.sh
